# install/nestor.service
[Unit]
Description=Nestor Server
After=network.target

[Service]
Type=simple
User=nestor
WorkingDirectory=/opt/nestor/app/server
ExecStart=/usr/bin/node dist/index.js
Environment=NODE_ENV=production
Environment=NESTOR_PORT=3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
